# dist/ubuntu/20.04/inc/etc/systemd/system/metrist-ebpf-agent.service
[Unit]
Description=Metrist.io Orchestrator eBPF plugin

[Service]
ExecStart=/usr/bin/metrist-ebpf-agent
EnvironmentFile=-/etc/default/metrist-ebpf-agent

[Install]
WantedBy=multi-user.target
